# systemd/.config/systemd/user/low_battery_norifier.timer
[Unit]
Description=Run the battery check script every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
