# systemd/user/low_battery_norifier.service
[Unit]
Description=Check battery level and send notifications

[Service]
Type=simple
ExecStart=/home/matvei/.config/scripts/check_battery.sh
